# scripts/com.financial-agent.automation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.financial-agent.automation</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/saumyatiwari/Documents/projects/financial-agent/scripts/run-automation.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>600</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/saumyatiwari/Documents/projects/financial-agent/logs/automation.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/saumyatiwari/Documents/projects/financial-agent/logs/automation.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>AUTOMATION_CRON_SECRET</key>
        <string>dev-secret-123</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
